# sync/crontab
*/15 * * * * python /app/config-sync.py
